# contrib/fuel_bootstrap/files/trusty/etc/rc.local
#!/bin/sh -e

fix-configs-on-startup || true

$(/usr/bin/let-rename)

