# contrib/fuel_bootstrap/files/trusty/lib/udev/rules.d/90-net-config-fix-on-rename.rules
# do not edit this file, it will be overwritten on update

ACTION=="move", SUBSYSTEM=="net", RUN+="fix-renamed-net-config $env{DEVPATH_OLD} $env{DEVPATH}"
