# contrib/fuel_bootstrap/files/trusty/etc/init/down-pxe-interface.conf
# down-pxe-interface
#
# The job downs the network interface used for PXE
# boot, what enableing udev to rename the interface  
# (in case of using a predictable interface naming
# schema). Otherwise udev will get "device busy"
# trying to rename the interface.
#
# The job is started on local-filesystems once at
# the beginning. The startup (event) is necessary
# in "start on " condition due to we want run the
# job only once.
#
# The pxe-interfaces-ready event is emitted,
# what is triggered the job 
# /etc/init/interfaces-are-ready.conf
# and in its turn enabling continue the job
# /etc/init/hold-interface.conf
#
# start on startup (event) is necessary due
# to we want run the job only once

start on startup and local-filesystems

emits pxe-interface-ready

task

pre-start script
# for debug uncomment following two lines
# exec >>/dev/.initramfs/pxe-interface-down.log 2>&1
# set -x
# Need the folder /run/network for ifdown 
   mkdir -p /run/network
end script

script 
# for debug uncomment following two lines
 exec >>/dev/.initramfs/pxe-interface-down.log 2>&1
 set -x

       # Lets down the interface which has been used for PXE and still up
       for DEV in /sys/class/net/* ; do
           INTERFACE=${DEV##*/}
           if [ "${INTERFACE}" != "lo" ]; then 
                 state="/sys/class/net/${INTERFACE}/operstate"
               if [ -f "${state}" ] && [ $(cat "${state}") = "up" ]; then 
                   ifdown --force  ${INTERFACE}
                   # trigger udev to run the rules and rename the interface
                   udevadm trigger --action=add --subsystem-match=net
                   udevadm settle
                   # interface-security has probably started, stop it
                   # initctl stop network-interface-security INTERFACE="${INTERFACE}"
               fi
           fi
       done
end script

post-start script 
# for debug uncomment following two lines
 exec >>/dev/.initramfs/pxe-interface-down.log 2>&1
 set -x
     # whethere the interface was down or not emit the event 
     initctl emit --no-wait pxe-interface-ready
end script

