# contrib/fuel_bootstrap/files/trusty/etc/init/hold-network-interface.conf
# hold-network-interface 
# 
# The dummy JOB to hold network interfaces until the
# network interface used for PXE will be put down 
# (if any), enabling the udev to change the iterface
# name according the rules (if any udev rule is to
# applied)
# 


start on starting network-interface-security and 
         started interfaces-are-ready

instance hold-network-interface${INTERFACE:+/}${INTERFACE:-}

task

exec :
