# contrib/fuel_bootstrap/files/trusty/etc/init/interfaces-are-ready.conf
# interfaces-are-ready
#
# The JOB is a trigger and the indicator, that
# we have got the event about interfaces' rediness.
# The interface used for PXE already is down and
# udev can rename it (if requested)
#

start on network-interfaces-ready

task

exec :
