# contrib/fuel_bootstrap/files/trusty/etc/init/interfaces-are-ready.conf
# interfaces-are-ready
#
# The JOB is a trigger and indicating, that
# we have got the event about interfaces' rediness.
# The interface used for PXE already is down and
# udev can rename it (if requested)
#
# It's not a task, so it will be at running state
# (task passes all states to the stopped)
# 

start on startup and pxe-interface-ready


pre-start script
 exec >>/dev/.initramfs/.log 2>&1
 set -x
 echo  "pre-start script"
end script

post-start script
 exec >>/dev/.initramfs/interfaces-are-ready.log 2>&1
 set -x
 echo  "post-start script"
end script

script 
 exec >>/dev/.initramfs/interfaces-are-ready.log 2>&1
 set -x
 echo  "script"
end script

pre-stop script
 exec >>/dev/.initramfs/interfaces-are-ready.log 2>&1
 set -x
 echo  "pre-stop script"
end script

post-stop script
 exec >>/dev/.initramfs/interfaces-are-ready.log 2>&1
 set -x
 echo  "post-stop script"
end script
