# contrib/fuel_bootstrap/files/trusty/etc/init/pxe-interface-down.conf
# pxe-interface-down
#
# The job downs the network interface used for PXE
# boot, what enables udev to rename the interface  
# (in case of using a predictable interface naming
# schema. Otherwise udev will get "device busy"
# trying to rename the interface according the rules.
#
# The job started on local-filesystems once at the
# beginning. 
# It emits the network-interfaces-ready event,
# which are triggered the job 
# /etc/init/interfaces-are-ready.conf
# and in its turn the enable continue the job
# /etc/init/hold-interface.conf
#

start on local-filesystems

emit network-interfaces-ready

task

script 
# for debug uncomment following two lines
 exec >>/dev/.initramfs/pxe-interface-down.log 2>&1
 set -x

       # Lets down the interface which has been used for PXE and still up
       for DEV in /sys/class/net/* ; do
           INTERFACE=${DEV##*/}
           if [ $(cat "/sys/class/net/${INTERFACE}/operstate")=="up" ]; then
               if [ "lo"!="${INTERFACE}" ]; then 
                   ifdown --force  ${INTERFACE}
                   # trigger udev to run the rules and rename the interface
                   udevadm trigger --action=add --subsystem-match=net
                   udevadm settle
               fi
           fi
       done
 
       # whethere the interface was down or not emit the event 
       initctl emit --no-wait network-interfaces-ready

end script

